# supervisor/siyu-worker.ini
[program:siyu-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /www/wwwroot/live.grs333.com/artisan queue:work redis --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/www/wwwroot/live.grs333.com/storage/logs/worker.log
